# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/nginx -g 'daemon on; master_process on;'